# udev-rules/80.power.rules
ACTION=="change", SUBSYSTEM=="power_supply", RUN+="/usr/bin/slimbookintelcontroller-pkexec"
